# platform/roles/gitlab-backup/templates/gitlab-backup.timer.j2
[Unit]
Description=Schedule GitLab backup pipeline

[Timer]
OnCalendar={{ systemd_timer_calendar }}
Persistent=true
AccuracySec=1m
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target

